# wireguard_client/rootfs/etc/services.d/wireguard_client/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: TorGuard WireGuard Client
# Runs TorGuard WireGuard Client
# ==============================================================================
declare interface

s6-svc -O /var/run/s6/services/wireguard_client

bashio::log.info "Starting TorGuard WireGuard Client..."

# This is alpha software. We need to set this to instruct
# WireGuard we are OK to go.
export WG_I_PREFER_BUGGY_USERSPACE_TO_POLISHED_KMOD=1

# Get the interface
interface="wg0"

# Run TorGuard WireGuard
exec wg-quick up "${interface}"
